# crates/cuda/src/cuslop.service.template
[Unit]
Description=cuslop-server-{%}

[Service]
ExecStart=%h/.sp1/bin/cuslop-server
Type=simple
Restart=on-failure
Environment="CUDA_VISIBLE_DEVICES={%}"

[Install]
WantedBy=default.target